# service/getty-tty2/run
#!/bin/execlineb -P

/sbin/agetty --noissue tty2
