# service/.s6-svscan/finish
#!/bin/execlineb -S0

cd /
redirfd -w 2 /dev/console
fdmove -c 1 2

foreground { s6-svc -X -- /run/s6/service/s6-svscan-log }
emptyenv -c
wait -r -- { }

/etc/s6/init-shutdown ${@}
